# drivetemp-drive-host-1.service
[Unit]
Description=Send drive cage temperatures via MQTT
After = network.target

[Service]
WorkingDirectory=/opt/drivetemp
ExecStart=/opt/drivetemp/drivetemp-drive-host.py 1
Restart=always
RestartSec=1s

[Install]
WantedBy = multi-user.target